# doc/lockd2-cert-reload.service
[Unit]
Description=lockd2 TLS tanúsítvány újratöltés (certbot deploy hook)
After=lockd2.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl reload lockd2.service
